# deploy/production-backup/scems-production-backup.timer
[Unit]
Description=Daily SCEMS sc_production paired backup

[Timer]
OnCalendar=*-*-* 02:30:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m
Unit=scems-production-backup.service

[Install]
WantedBy=timers.target
