# deploy/production-backup/scems-production-backup.service
[Unit]
Description=SCEMS sc_production paired database and filestore backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/scems/production-backup.env
EnvironmentFile=/etc/scems/production-backup-runtime.env
ExecStart=/usr/bin/python3 /opt/ops/production_backup_restore.py backup
User=root
Group=root
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
